# .s6/etc/services.d/gglass/run
#!/bin/execlineb -P

# with-contenv
s6-envdir -fn -- /var/run/s6/container_environment

# Redirect stderr to stdout.
fdmove -c 2 1

cd /app

# TODO: Add support to skip this if redis wasn't going to start anyway
if { s6-svwait -U /var/run/s6/services/redis/ }

s6-setuidgid abc

/usr/bin/node /app/dist/server.js
#/usr/bin/npm start
